# monitor/backend/supervisor/backend.conf
[program:backend]
process_name=backend
priority=500
numprocs=1
user=<specify your user here>
directory=/path/to/main/directory/
command=python3.7 main.py
autorestart=true
