# servers/monitor/backend/supervisor/oracle.conf
[program:oracle]
process_name=oracle
priority=0
numprocs=1
user=<specify user here>
directory=/path/to/oracle
command=python3.7 oracle.py
autorestart=true
